# tests/funcional/assets/cron.d/sei

SHELL=/bin/bash

* * * * * root . /root/crond_env.sh; /usr/bin/php -c /etc/php.ini /opt/sei/scripts/AgendamentoTarefaSEI.php 2>&1 >> /var/log/sei/agendamento_sei.log
00 01 * * * root rm -rf /opt/sei/temp/*
